# extra/foreman-proxy-omaha-sync.cron
# Sync omaha releases once a day
0 22 * * * foreman-proxy source /opt/theforeman/tfm/enable; smart-proxy-omaha-sync >>/var/log/foreman-proxy/cron.log 2>&1
